# iosevka-generated.service
[Unit]
Description=Run iosevka-generate when config changes

[Service]
Restart=always

[Install]
WantedBy=default.target
